# src/service/Hermesd
#!/bin/sh
# chkconfig: 2345 20 80
# description:


. /etc/rc.d/init.d/functions

    
#NAME=/root/bot/Hermes
   

case "$1" in
start)
   nohup /usr/sbin/Hermes.sh 1>> /var/log/Hermes/Hermes.log 2>>/var/log/Hermes/Hermes.err &
   echo $!>/var/run/Hermes.sh.pid
   ;;
stop)
   pkill -P `cat /var/run/Hermes.sh.pid`
   rm /var/run/Hermes.sh.pid
   ;;
restart)
   $0 stop
   $0 start
   ;;
status)
   if [ -e /var/run/Hermes.sh.pid ]; then
      echo Hermes.sh is running, pid=`cat /var/run/Hermes.sh.pid`
   else
      echo Hermes.sh is not running
      exit 1
   fi
   ;;
*)
   echo "Usage: $0 {start|stop|status|restart}"
esac

exit 0 
